# supporting_files/supervisord-jetty.conf
[program:jetty]
command=/start-jetty.sh
numprocs=1
autostart=true
autorestart=true