# supporting_files/supervisord-jetty.conf
[program:jetty]
directory=$JETTY_BASE
command= java -jar $JETTY_HOME/start.jar
user=jetty
numprocs=1
autostart=true
autorestart=true
startsecs=30
